# rpi/openapi-server.service
# openapi_server.service
[Unit]
Description=OpenAPI Server
After=network.target

[Service]
ExecStart=/usr/bin/python3 -m openapi_server
Restart=always
User=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target